# templates/getmail.service.erb
[Unit]
Description=Getmail for <%= @imap_username %>

[Service]
User=getmail
ExecStart=/usr/bin/getmail --getmaildir=/etc/getmail/<%= @imap_username %> --idle INBOX
Restart=always

[Install]
WantedBy=multi-user.target
